# deploy/systemd/board-wired-network.service
[Unit]
Description=Configure SL1680 wired dog link
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/aicps/robot-dog-voice-interface
ExecStart=/bin/sh /home/aicps/robot-dog-voice-interface/scripts/board_set_wired_ip.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
